# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:acestream]
command=acestreamengine --client-console --bind-all --login test --password test --live-buffer %(ENV_LIVEBUFFER)s --max-connections %(ENV_MAXCONNECTIONS)s --max-peers %(ENV_MAXPEERS)s --download-limit %(ENV_DLIMIT)s --upload-limit %(ENV_ULIMIT)s --log-stdout --port %(ENV_PUBPORT)s
directory = /home/tv/

[program:acehttp]
command=/home/tv/aceproxy-master/acehttp.py
user = tv
directory = /home/tv/
autorestart = true

[program:vlc]
command=vlc -I telnet --clock-jitter 0 --sout-mux-caching 500 --network-caching %(ENV_NETCACHING)s --telnet-password admin --telnet-host 127.0.0.1
environment=VLC_PLUGIN_PATH="/usr/lib/vlc/plugins/"
user = tv
directory = /home/tv/
